# scripts/com.spelab.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spelab.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>-m</string>
        <string>server.app</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/fengchangrui/Desktop/cursor/SPELAB</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/fengchangrui/Library/Logs/SPELAB/server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fengchangrui/Library/Logs/SPELAB/server.log</string>
</dict>
</plist>
